# files/minecraft.service
[Unit]
Description=Start [SERVER] minecraft server

[Service]
Type=oneshot
ExecStart=[SERVER-PATH]/start.sh
